# deployment/files/app.service.j2
[Unit]
Description={{ app_display_name }} Flask Application
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment="PATH={{ venv_dir }}/bin"
Environment="FLASK_APP=app"
Environment="APP_SERVICE_NAME={{ service_name }}"

# Load environment variables from .env file
EnvironmentFile={{ app_dir }}/.env

# Run gunicorn with proper logging
ExecStart={{ venv_dir }}/bin/gunicorn \
    --bind {{ gunicorn_bind }} \
    --workers {{ gunicorn_workers }} \
    --threads 2 \
    --timeout {{ gunicorn_timeout }} \
    --access-logfile {{ log_dir }}/access.log \
    --error-logfile {{ log_dir }}/error.log \
    --log-level warning \
    "app:create_app()"

# Restart policy
Restart=always
RestartSec=10

# ============================================================================
# Security Hardening
# ============================================================================
# Prevents privilege escalation
NoNewPrivileges=true

# Private /tmp directory (isolated from other services)
PrivateTmp=true

# Restrict system calls to necessary ones only
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount

# Prevent access to /home directories (except its own working directory)
ProtectHome=read-only

# Make /usr, /boot, /etc read-only
ProtectSystem=strict

# Allow writes only to specific directories
ReadWritePaths={{ log_dir }}
ReadWritePaths={{ app_dir }}/instance

# Restrict access to kernel logs and modules
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true

# Prevent loading kernel modules
LockPersonality=true

# Restrict access to real-time clock
ProtectClock=true

# Private /dev (only necessary device nodes)
PrivateDevices=true

# Restrict network namespace
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Prevent access to other users' processes
ProtectProc=invisible

# No access to hostname changes
ProtectHostname=true

# Prevent gaining new privileges
RestrictSUIDSGID=true

# Memory execution protection
MemoryDenyWriteExecute=false

# Restrict namespace usage
RestrictNamespaces=true

# Remove capabilities (app doesn't need special privileges)
CapabilityBoundingSet=
AmbientCapabilities=

# Prevent access to control groups
ProtectControlGroups=true

# Umask for created files
UMask=0027

# ============================================================================
# Logging
# ============================================================================
StandardOutput=append:{{ log_dir }}/app.log
StandardError=append:{{ log_dir }}/app.log

[Install]
WantedBy=multi-user.target

