# deployment/templates/app.service.j2
# ============================================================================
# SYSTEMD SERVICE UNIT - APPLICATION WITH SECURITY HARDENING
# ============================================================================
# Template: deployment/templates/app.service.j2
# Purpose:
#   - Manages application lifecycle (start, stop, restart, enable on boot)
#   - Runs as unprivileged user with minimal permissions
#   - Enforces security constraints via systemd sandboxing
#   - Logs to journalctl with structured output
#
# Generated by: setup.yml playbook during deployment
# Location on server: /etc/systemd/system/{{ app_name }}.service
#
# Key Features:
#   ✓ Comprehensive security hardening (AppArmor-like restrictions)
#   ✓ Restricted system call filter (only necessary syscalls allowed)
#   ✓ Read-only root filesystem (only specific dirs writable)
#   ✓ No privilege escalation possible
#   ✓ Limited device access (no raw sockets)
#   ✓ Automatic restart on failure
#   ✓ Gunicorn with gthread workers
#   ✓ Environment variables for AWS Secrets Manager
#
# Usage:
#   systemctl start {{ app_name }}      # Start the service
#   systemctl stop {{ app_name }}       # Stop the service
#   systemctl restart {{ app_name }}    # Restart the service
#   systemctl status {{ app_name }}     # View status and recent logs
#   journalctl -u {{ app_name }} -f     # Follow logs in real-time
#
# ⚠️  DO NOT manually edit - will be overwritten on redeploy
#     Edit group_vars/vault.yml and rerun setup.yml instead
# ============================================================================

# ============================================================================
# SERVICE CONFIGURATION
# ============================================================================
# Type=simple: Gunicorn process runs in foreground (systemd manages it)
# User/Group: Run as non-root app_user (security best practice)
# WorkingDirectory: Set before execution (for relative paths)
# ============================================================================
[Service]
Type=simple
User={{ app_user }}
Group={{ app_name }}
WorkingDirectory={{ app_dir }}
Environment="PATH={{ venv_dir }}/bin"
Environment="FLASK_APP=app"
Environment="APP_SERVICE_NAME={{ service_name }}"

# ============================================================================
# AWS CONFIGURATION FOR SECRETS MANAGER
# ============================================================================
# These environment variables tell the app to read secrets from AWS
# The app config.py uses these to fetch secrets from Secrets Manager
#
# SECRET_NAME: The AWS Secrets Manager secret name
#              Format: {app_name}/production
# AWS_REGION: The AWS region where the secret is stored
# ============================================================================
Environment="SECRET_NAME={{ secret_name }}"
Environment="AWS_REGION={{ aws_region }}"

# ============================================================================
# APPLICATION CONFIGURATION
# ============================================================================
# Non-secret infrastructure config the app reads via os.environ.get()
# All actual secrets (keys, passwords, tokens) come from Secrets Manager
# ============================================================================
Environment="COMIC_IMAGE_PATH={{ app_dir }}/instance/item_images"
Environment="CLOUDWATCH_NAMESPACE={{ app_display_name }}"

# ============================================================================
# APPLICATION STARTUP WITH GUNICORN
# ============================================================================
# Gunicorn: WSGI HTTP server (multi-worker, multi-threaded)
#
# Configuration:
#   --workers: Number of processes (CPU-bound multiplier)
#   --threads: Threads per worker (I/O-bound multiplier)
#   --worker-class gthread: Thread-based workers (good for I/O)
#   --timeout: Kill worker if request takes longer than this
#   --bind: Where to listen (127.0.0.1:8000 = localhost only)
#   --access-logfile: Track all HTTP requests
#   --error-logfile: Application exceptions and errors
#   --log-level warning: Only log warnings and above
#   --capture-output: Redirect print() and logging to stdout
#
# Load app from: runapp module, app function
# ============================================================================
# Run gunicorn with proper logging
ExecStart={{ venv_dir }}/bin/gunicorn \
    --bind {{ gunicorn_bind }} \
    --workers {{ gunicorn_workers }} \
    --threads 2 \
    --timeout {{ gunicorn_timeout }} \
    --worker-tmp-dir /tmp \
    --control-socket /tmp/gunicorn.ctl \
    --access-logfile {{ log_dir }}/access.log \
    --error-logfile {{ log_dir }}/error.log \
    --log-level warning \
    "app:create_app('production')"

# ============================================================================
# RESTART POLICY
# ============================================================================
# Restart=always: Automatically restart on ANY exit (success or failure)
# RestartSec=10: Wait 10 seconds before restarting (prevents restart loop)
# ============================================================================
# Restart policy
Restart=always
RestartSec=10

# ============================================================================
# COMPREHENSIVE SECURITY HARDENING
# ============================================================================
# These settings restrict the process using systemd sandboxing
# Goal: Even if app is compromised, damage is limited
#
# Defense Layers:
#   1. Privilege restrictions: Can't gain new capabilities
#   2. Filesystem isolation: Read-only root, specific write paths
#   3. System call filtering: Only allow necessary syscalls
#   4. Device isolation: No raw device access
#   5. Namespace restrictions: Limited visibility into system
#
# Performance: <1% overhead, security gain is significant
# ============================================================================

# ============================================================================
# PRIVILEGE RESTRICTIONS
# ============================================================================
# Prevents privilege escalation
NoNewPrivileges=true

# Private /tmp directory (isolated from other services)
PrivateTmp=true

# ============================================================================
# SYSTEM CALL FILTERING (SecComp)
# ============================================================================
# Whitelist approach: Only allow specific system calls
#
# @system-service: Pre-defined list for typical services
#                 (fork, execve, mmap, read, write, etc.)
# ~@privileged: REMOVE privileged calls (mount, reboot, etc.)
# ~@resources: REMOVE resource-intensive calls
# ~@mount: REMOVE filesystem mounting calls
#
# Blocked syscalls: adjtimex, clock_adjtime, mbind, etc.
# If app tries blocked syscall: Process killed immediately
# ============================================================================
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount

# ============================================================================
# FILESYSTEM RESTRICTIONS
# ============================================================================
# ProtectHome=true: Completely block access to /home, /root, /run/user
#                  (app lives under /opt, does not need /home)
#
# ProtectSystem=strict: Root filesystem (/,/usr,/boot,/etc) read-only
#                      (prevents modifying system)
#
# ReadWritePaths: WHITELIST specific directories that must be writable
#                - Log directory
#                - Instance directory (for cached data)
#
# Default: Everything else is read-only
# ============================================================================
# Block all access to /home (app lives under /opt)
ProtectHome=true

# Make /usr, /boot, /etc read-only
ProtectSystem=strict

# Allow writes only to specific directories
ReadWritePaths={{ log_dir }}
ReadWritePaths={{ app_dir }}/instance

# ============================================================================
# KERNEL PROTECTION
# ============================================================================
# Restrict access to kernel logs, modules, and parameters
# Prevents tampering with system via kernel interface
# ============================================================================
# Restrict access to kernel logs and modules
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true

# Prevent loading kernel modules
LockPersonality=true

# Restrict access to real-time clock
ProtectClock=true

# ============================================================================
# DEVICE ISOLATION
# ============================================================================
# PrivateDevices=true: Only /dev/null, /dev/zero, /dev/full, /dev/random
#                     (no raw devices like /dev/sda, /dev/kmem)
# ============================================================================
# Private /dev (only necessary device nodes)
PrivateDevices=true

# ============================================================================
# NETWORK ISOLATION
# ============================================================================
# Restrict to standard network families only
# AF_UNIX: Local sockets
# AF_INET: IPv4
# AF_INET6: IPv6
# ============================================================================
# Restrict network namespace
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# ============================================================================
# PROCESS AND NAMESPACE ISOLATION
# ============================================================================
# ProtectProc=invisible: Other processes are invisible
# ProtectHostname=true: Can't change system hostname
# ============================================================================
# Prevent access to other users' processes
ProtectProc=invisible

# No access to hostname changes
ProtectHostname=true

# ============================================================================
# PRIVILEGE AND CAPABILITY RESTRICTIONS
# ============================================================================
# RestrictSUIDSGID=true: Can't execute SUID/SGID binaries
# MemoryDenyWriteExecute=false: Allow JIT compilation (Python doesn't need)
# RestrictNamespaces=true: Can't unshare namespaces
# ============================================================================
# Prevent gaining new privileges
RestrictSUIDSGID=true

# Memory execution protection
MemoryDenyWriteExecute=false

# Restrict namespace usage
RestrictNamespaces=true

# Remove capabilities (app doesn't need special privileges)
# CapabilityBoundingSet empty = remove ALL capabilities
# AmbientCapabilities empty = don't add any
CapabilityBoundingSet=
AmbientCapabilities=

# Prevent access to control groups
ProtectControlGroups=true

# ============================================================================
# FILE CREATION MASK
# ============================================================================
# UMask=0027: Files created with permissions 640 (rw-r-----)
#            Prevents world-readable sensitive files
# ============================================================================
# Umask for created files
UMask=0027

# ============================================================================
# LOGGING CONFIGURATION
# ============================================================================
# StandardOutput: Application logs (stdout) written to append mode
# StandardError: Application errors (stderr) written to append mode
#
# View logs:
#   journalctl -u {{ app_name }} -f          # Real-time follow
#   journalctl -u {{ app_name }} -n 100      # Last 100 lines
#   journalctl -u {{ app_name }} --since "1 hour ago"
#   tail -f {{ log_dir }}/app.log            # Direct file view
#
# Log rotation: Handled by logrotate (see deployment/templates/logrotate.conf.j2)
# ============================================================================
# Logging
StandardOutput=append:{{ log_dir }}/app.log
StandardError=append:{{ log_dir }}/app.log

[Install]
WantedBy=multi-user.target

