# deployment/templates/supervisor.conf.j2
# ============================================================================
# SUPERVISOR PROCESS MANAGER CONFIGURATION
# ============================================================================
# Template: deployment/templates/supervisor.conf.j2
# Purpose:
#   - Alternative to systemd for process management
#   - Keeps application running if it crashes
#   - Easy start/stop/restart from supervisorctl CLI
#
# Note: This is an alternative to systemd-with-validation.service.j2
#       Choose ONE for your deployment (systemd is recommended on Ubuntu)
#
# Generated by: setup.yml playbook (if supervisor is preferred)
# Location on server: /etc/supervisor/conf.d/{{ app_name }}.conf
#
# Usage:
#   supervisorctl status {{ service_name }}   # Check status
#   supervisorctl start {{ service_name }}    # Start process
#   supervisorctl stop {{ service_name }}     # Stop process
#   supervisorctl restart {{ service_name }}  # Restart process
#   supervisorctl tail {{ service_name }}     # Follow logs
# ============================================================================

[program:{{ service_name }}]
# ============================================================================
# APPLICATION STARTUP COMMAND
# ============================================================================
# Gunicorn configuration:
#   -w: Number of worker processes
#   -b: Bind address and port (127.0.0.1:8000 = localhost only)
#   --timeout: Kill worker if request takes longer
#   --access-logfile: Track all HTTP requests
#   --error-logfile: Gunicorn process errors (worker crashes, timeouts)
#   --log-level: Verbosity (debug, info, warning, error)
#
#   NOTE: Do NOT use --capture-output. Flask writes its own app.log and error.log
#         via RotatingFileHandler. --capture-output would duplicate INFO into
#         Gunicorn's error.log.
# ============================================================================
command={{ venv_dir }}/bin/gunicorn -w {{ gunicorn_workers }} -b 127.0.0.1:8000 --timeout {{ gunicorn_timeout }} --access-logfile {{ log_dir }}/access.log --error-logfile {{ log_dir }}/error.log --log-level warning "app:create_app('production')"

# Working directory for the process
directory={{ app_dir }}

# Run as unprivileged user (security best practice)
user={{ app_user }}

# Auto-start when supervisor starts
autostart=true

# Auto-restart if process exits
autorestart=true

# Kill all child processes when stopping
stopasgroup=true
killasgroup=true

# Log file locations
stderr_logfile={{ log_dir }}/app.log
stdout_logfile={{ log_dir }}/app.log

# Environment variables (PATH must include venv, SECRET_NAME/AWS_REGION for Secrets Manager)
environment=PATH="{{ venv_dir }}/bin",SECRET_NAME="{{ secret_name }}",AWS_REGION="{{ aws_region }}",APP_SERVICE_NAME="{{ service_name }}",APP_NAME="{{ app_name }}"
