# deployment/templates/systemd-with-validation.service.j2
# ============================================================================
# SYSTEMD SERVICE UNIT - APPLICATION WITH PRE-START VALIDATION
# ============================================================================
# Template: deployment/templates/systemd-with-validation.service.j2
# Purpose:
#   - Manages application lifecycle (start, stop, restart, enable on boot)
#   - Runs validation BEFORE attempting to start
#   - Prevents startup if configuration invalid
#   - Logs application output to journalctl
#
# Generated by: setup.yml playbook during deployment
# Location on server: /etc/systemd/system/{{ app_name }}.service
#
# Key Features:
#   ✓ Pre-start CSV schema validation (catches config errors early)
#   ✓ Automatic restart on failure (RestartSec=10)
#   ✓ Runs as unprivileged application user (security hardening)
#   ✓ Systemd socket activation ready
#   ✓ Resource limits (file descriptor limit: 65536)
#   ✓ Gunicorn with gthread workers (multi-threading for I/O bound tasks)
#   ✓ Environment variables for AWS Secrets Manager
#   ✓ Separate error and access logs
#
# Usage:
#   systemctl start {{ app_name }}      # Start the service
#   systemctl stop {{ app_name }}       # Stop the service
#   systemctl restart {{ app_name }}    # Restart the service
#   systemctl status {{ app_name }}     # View status and recent logs
#   journalctl -u {{ app_name }} -f     # Follow logs in real-time
#
# ⚠️  DO NOT manually edit - will be overwritten on redeploy
#     Edit group_vars/vault.yml and rerun setup.yml instead
#
# Variables from vault.yml:
#   - app_display_name: Friendly name for service description
#   - app_user: Non-privileged user to run as
#   - app_dir: Application directory (on EBS volume)
#   - venv_dir: Python virtual environment
#   - log_dir: Application logs directory
#   - secret_name, aws_region: AWS configuration
# ============================================================================

# ============================================================================
# SERVICE CONFIGURATION
# ============================================================================
# Type=simple: Gunicorn process runs in foreground (systemd manages it)
# User/Group: Run as non-root app_user (security best practice)
# WorkingDirectory: Set before execution (for relative paths)
# ============================================================================
[Service]
Type=simple
User={{ app_user }}
Group={{ app_name }}
WorkingDirectory={{ app_dir }}

# ============================================================================
# ENVIRONMENT CONFIGURATION
# ============================================================================
# PATH: Ensure Python venv is first in path (for correct pip/python versions)
# PYTHONUNBUFFERED=1: Disable Python buffering (see logs in real-time)
# SECRET_NAME: Points to AWS Secrets Manager secret
# AWS_REGION: Required for AWS API calls
# ============================================================================
# Environment
Environment="PATH={{ venv_dir }}/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="SECRET_NAME={{ secret_name }}"
Environment="AWS_REGION={{ aws_region }}"
Environment="APP_SERVICE_NAME={{ service_name }}"
Environment="COMIC_IMAGE_PATH={{ app_dir }}/instance/item_images"
Environment="CLOUDWATCH_NAMESPACE={{ app_display_name }}"

# ============================================================================
# PRE-START VALIDATION
# ============================================================================
# ExecStartPre runs BEFORE ExecStart and must succeed (exit code 0)
# If validation fails, the service will not start
#
# This validates:
#   1. CSV data files are in correct schema
#   2. Configuration is loadable
#   3. Dependencies are available
#
# Why validate before start?
#   - Prevents silent failures
#   - Catches configuration errors early
#   - Clear error messages in systemd logs
#   - Faster debugging vs. app crashing after start
# ============================================================================
# Pre-start validation - PREVENTS startup if schema invalid
ExecStartPre=/bin/bash -c 'echo "Validating CSV schema..."'
ExecStartPre={{ venv_dir }}/bin/python {{ app_dir }}/app/scripts/util_validate_csv_schema.py --quiet
ExecStartPre=/bin/bash -c 'echo "✓ CSV schema validation passed"'

# ============================================================================
# APPLICATION STARTUP WITH GUNICORN
# ============================================================================
# Gunicorn: WSGI HTTP server (multi-worker, multi-threaded)
#
# Configuration:
#   --workers 2: Run 2 processes (handles 2 concurrent long-running requests)
#   --threads 4: Each worker gets 4 threads (handles 8 concurrent total)
#   --worker-class gthread: Thread-based workers (good for I/O bound apps)
#   --timeout 120: Kill worker if request takes > 120 seconds
#   --bind 127.0.0.1:8000: Listen on localhost (nginx proxies to here)
#
# Logs:
#   - access log: Track all requests (GET /api/comics, 200, 45ms)
#   - error log: Errors and exceptions from app
#
# Tuning:
#   - Increase workers if CPU-bound (make sure you have enough RAM)
#   - Increase threads if I/O-bound (DB queries, API calls)
#   - Increase timeout if long-running tasks
# ============================================================================
# Start application
# Using 2 workers to minimize startup logging duplication
# Each worker gets 4 threads for concurrency (2 workers * 4 threads = 8 concurrent requests)
ExecStart={{ venv_dir }}/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --threads 4 \
    --worker-class gthread \
    --timeout 120 \
    --access-logfile {{ log_dir }}/access.log \
    --error-logfile {{ log_dir }}/error.log \
    --log-level info \
    --capture-output \
    runapp:app

# ============================================================================
# RESTART POLICY
# ============================================================================
# Restart=always: Automatically restart on ANY exit (success or failure)
# RestartSec=10: Wait 10 seconds before restarting (prevents restart loop)
#
# Note: If validation fails, RestartSec gives time to investigate logs
#       Then systemd automatically retries (good for transient errors)
# ============================================================================
# Restart policy
Restart=always
RestartSec=10

# ============================================================================
# SECURITY HARDENING
# ============================================================================
# These settings restrict what the process can do (defense in depth)
#
# NoNewPrivileges=true: Prevents gaining privileges via setuid binaries
# PrivateTmp=true: Gives isolated /tmp (other processes can't see files)
# ============================================================================
# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# ============================================================================
# LOGGING CONFIGURATION
# ============================================================================
# StandardOutput: App logs written to /var/log/...
# StandardError: App errors written to /var/log/...
#
# View logs:
#   journalctl -u {{ app_name }} -f          # Real-time follow
#   journalctl -u {{ app_name }} -n 50       # Last 50 lines
#   journalctl -u {{ app_name }} --since "2 hours ago"
#
# Logs are also stored in {{ log_dir }}/app.log and {{ log_dir }}/error.log
# ============================================================================
# Logging
StandardOutput=append:{{ log_dir }}/app.log
StandardError=append:{{ log_dir }}/error.log

# ============================================================================
# RESOURCE LIMITS
# ============================================================================
# LimitNOFILE=65536: Maximum open file descriptors
#
# Why this matters:
#   - S3 uploads/downloads need multiple connections
#   - Default is usually 1024 (too low for concurrent users)
#   - 65536 is typical for servers handling many connections
# ============================================================================
# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
